# src/tools/anope.service.in
[Unit]
After=network.target
Description=Anope IRC Services
Documentation=https://wiki.anope.org/
After=network-online.target
Wants=network-online.target

[Service]
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=@BIN_DIR@/@PROGRAM_NAME@ --nofork
Restart=on-failure
Type=simple
WorkingDirectory=@CMAKE_INSTALL_PREFIX@

[Install]
WantedBy=multi-user.target

